# metadata-generator/scheduling/plex-metadata-generator.timer
[Unit]
Description=Plex Metadata Generator Daily Schedule
Documentation=https://github.com/roto31/plex-metadata-generator
Requires=plex-metadata-generator.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00

# Run on boot if missed
Persistent=true

# Only run if system is idle (no high load)
ConditionACPower=true

# Random delay to prevent multiple tasks at exact same time
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
