# metadata-generator/scheduling/plex-metadata-generator.service
[Unit]
Description=Plex Metadata NFO Generator with Tunarr/TVDb/TMDb Integration
Documentation=https://github.com/roto31/plex-metadata-generator
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=plex
Group=plex
ExecStart=/usr/local/bin/plex-metadata-generator --config /etc/plex-metadata-generator.conf --no-prompts

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/mnt/media/TV /var/cache/plex-metadata /var/log

# Resource limits
CPUQuota=50%
MemoryLimit=1G
TasksMax=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plex-metadata-gen

[Install]
WantedBy=multi-user.target
