# metadata-generator/scheduling/plex-metadata-generator-windows.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Plex Metadata Generator — Windows Task Scheduler

  Import via Task Scheduler UI:
    1. Open Task Scheduler  (taskschd.msc)
    2. Action → Import Task...
    3. Select this file
    4. Update the <UserId> and script path to match your system
    5. Click OK

  Or via PowerShell (run as Administrator):
    Register-ScheduledTask -Xml (Get-Content .\plex-metadata-generator-windows.xml | Out-String) `
      -TaskName "PlexMetadataGenerator" -Force
-->
<Task version="1.4"
  xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>
    <Date>2026-06-17T00:00:00</Date>
    <Author>roto31</Author>
    <Description>Runs Plex Metadata Generator daily at 2 AM to update NFO metadata files and refresh Plex libraries.</Description>
    <URI>\PlexMetadataGenerator</URI>
  </RegistrationInfo>

  <Triggers>
    <!-- Daily at 2:00 AM -->
    <CalendarTrigger>
      <StartBoundary>2026-06-17T02:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <!-- Replace with your Windows username -->
      <UserId>YOUR_WINDOWS_USERNAME</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <!-- Start within 1 hour if the scheduled time was missed -->
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT2H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <!-- Update this path to match where Python is installed -->
      <Command>C:\Python312\python.exe</Command>
      <!-- Update this path to where you placed the script -->
      <Arguments>"C:\Program Files\PlexMetadataGenerator\plex_metadata_generator_extended.py" --config "C:\ProgramData\PlexMetadataGenerator\plex-metadata-generator.conf" --media-type all</Arguments>
      <WorkingDirectory>C:\Program Files\PlexMetadataGenerator</WorkingDirectory>
    </Exec>
  </Actions>

</Task>
